# flask-app-v5/etc/profile.d/sh.local
if [ -f /opt/app-root/bin/activate ]; then
    source /opt/app-root/bin/activate
fi
